# backend/com.clubhub.api.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clubhub.api</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/devinmcnair/Desktop/ClubHub/backend/clubhub_api.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/devinmcnair/Desktop/ClubHub/backend</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CLUBHUB_DEBUG</key>
    <string>0</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/devinmcnair/Desktop/ClubHub/backend/logs/launchd.api.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/devinmcnair/Desktop/ClubHub/backend/logs/launchd.api.err.log</string>
</dict>
</plist>
